# swarm/build/dist/etc/profile.d/p4cli.sh
PATH=/opt/perforce/usr/bin:$PATH; export PATH
